# telme10.socket
[Unit]
Description=realraum birthday telnet server socket

[Socket]
ListenStream=23

[Install]
WantedBy=sockets.target
